# root/etc/s6-overlay/s6-rc.d/svc-zigbee2mqtt/run
#!/usr/bin/with-contenv bash

cd /app/zigbee2mqtt || exit 1

echo "Starting Zigbee2MQTT."
exec s6-setuidgid abc \
	node index.js
