# root/etc/s6-overlay/s6-rc.d/svc-mosquitto/run
#!/usr/bin/with-contenv bash

echo "Starting Mosquitto."
exec s6-setuidgid abc \
	mosquitto \
	-c /config/mosquitto/mosquitto.conf
