# root/etc/services.d/zigbee2mqtt/run
#!/usr/bin/with-contenv bash

cd /app/zigbee2mqtt || exit 1

echo "Starting Zigbee2MQTT."
redirect_cmd exec s6-setuidgid abc \
	node index.js
