# root/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bash

echo "Starting Mosquitto."
redirect_cmd exec s6-setuidgid abc \
	mosquitto \
	-c /config/mosquitto/mosquitto.conf
